# net-nds/pamldapd/files/pamldapd.initd
#!/sbin/openrc-run
name="pamldapd daemon"
description="Simple LDAP server, uses PAM as backend"
command=/usr/bin/pamldapd
command_args="-c /etc/pamldapd.json -l /var/log/pamldapd/pamldapd.log"
command_background=true
command_user="pamldapd:pamldapd"
output_log="/var/log/pamldapd/stdout.log"
error_log="/var/log/pamldapd/stdout.log"
pidfile=/var/run/pamldapd.pid
